# root/etc/services.d/flexget/run
#!/usr/bin/with-contenv bash

exec \
       s6-setuidgid abc /usr/bin/flexget -c /flexcfg/config.yml \
       -l /flexcfg/flexget.log -L info daemon start -d --autoreload-config
